# scripts/RepImmo/repimmo.crontab
* * * * * /home/REPIMMO/prg/scripts/RepImmoRequete.sh > /log/REPIMMO/oracle/RepImmoRequete.log
